# makepkg/pkg/tea-utils/teaOS/utils/sys-init/sys-init.service
[Unit]
Description=teaOS sys-init script

[Service]
Type=simple
ExecStart=/bin/sh /teaOS/utils/sys-init/run

[Install]
WantedBy=default.target
